# pacman-to-cache.timer
[Unit]
PartOf=pacman-to-cache.service

[Timer]
OnActiveSec=10m

[Install]
WantedBy=multi-user.target
